# fleet/build-io-base/build-io-base.service
[Unit]
Description=Build iobase Docker image
After=docker-registry.service

[Service]
Type=oneshot
ExecStart=/bin/sh /opt/data/fleet/build-io-base/start.sh
RemainAfterExit=yes
