# fleet/arkenwatch/arkenwatch.service
[Unit]
Description=arkenwatch service

[Service]
ExecStart=/opt/data/fleet/arkenwatch/start.sh
ExecStop=/opt/data/fleet/arkenwatch/stop.sh
