# systemd/datadog/datadog.service
[Unit]
Description=datadog service

[Service]
ExecStart=/opt/data/systemd/datadog/start.sh
ExecStop=/opt/data/systemd/datadog/stop.sh
